# Downloads/saliency-20140505/packages/hacks/pixci
#!/bin/sh -e
### BEGIN INIT INFO
# Provides:          pixci
# Required-Start:    $local_fs gdm
# Required-Stop:
# Default-Start:     1 2 3 4 5
# Default-Stop:
# Short-Description: Start PIXCI driver
### END INIT INFO
#
# Start PIXCI driver

case "$1" in
  start)
        # insert the module:
	insmod /usr/local/xcap/drivers/x86_64/2.6.38-11-generic/pixci_x86_64.ko

        # create the device:
	node=`cat /proc/devices | /bin/grep PIXCI | awk '{ print $1 }'`
	rm -f /dev/pixci
	mknod /dev/pixci c $node 0
	chmod 777 /dev/pixci
	;;
    stop|restart|force-reload) exit 0 ;;
    *) echo "Usage: $0 {start|stop|restart|force-reload}" >&2; exit 1 ;;
esac
